# .infra/ansible/roles/jenkins/templates/jenkins-agent.service.j2
[Unit]
Description=Jenkins Agent

[Service]
User={{ user }}
WorkingDirectory=/home/cicd
ExecStart=/bin/bash {{ script_path }}
Restart=always

[Install]
WantedBy=multi-user.target
